# src/install/systemd/backup.service
[Unit]
Description=Backup Script

[Service]
Type=oneshot
ExecStart=/data/Backup/backup.sh
User=root
