# src/install/systemd/node_exporter.service
[Unit]
Description=Node Exporter

[Service]
#Type=forking
ExecStart=/sbin/node_exporter
Restart=always

[Install]
WantedBy=multi-user.target
