# install/at-package/at-3.1.23/atd.service.in
[Unit]
Description=Deferred execution scheduler
Documentation=man:atd(8)
After=remote-fs.target nss-user-lookup.target

[Service]
ExecStartPre=-find @atjobdir@ -type f -name "=*" -not -newercc /run/systemd -delete
ExecStart=@sbindir@/atd -f
IgnoreSIGPIPE=false
KillMode=process
Restart=on-failure

[Install]
WantedBy=multi-user.target
